# web/launchd/tradingagents-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradingagents.web</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/a1/TradingAgents/.venv/bin/uvicorn</string>
    <string>web.app:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8000</string>
    <string>--no-access-log</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- Fail-closed: set here (NOT in .env) so the server refuses to start
         ungated if .env is ever missing/renamed. The password itself lives
         in .env. -->
    <key>TRADINGAGENTS_WEB_REQUIRE_PASSWORD</key>
    <string>1</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/a1/TradingAgents</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/a1/.tradingagents/logs/web-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/a1/.tradingagents/logs/web-stderr.log</string>
</dict>
</plist>
